# files/default/hung_service.conf
# Hung_service Service

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

respawn
respawn limit 2 5

exec false

post-start script
    sleep 600
end script